# config/systemd/user/tmux-resurrect-cleanup.timer
[Unit]
Description=Weekly cleanup of old tmux-resurrect snapshots

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
